# indigo-core/common-rootfs-template/etc/sv/httpd/run
#!/bin/sh
exec httpd -f -h /usr/local/httpd 2>&1
