# indigo-gsm73xx/platform-rootfs-template/etc/sv/ofswd/finish
#!/bin/sh

# when stopping ofswd, restart inband-network-cfg 
source /etc/find-env
sv restart inband-network-cfg
