# indigo-gsm73xx/platform-rootfs-template/etc/inittab
::sysinit:/etc/init.d/rcS
::ctrlaltdel:/sbin/reboot
::shutdown:/bin/umount -a -r
::restart:/sbin/init
::respawn:/sbin/runsvdir-start
ttyS0::respawn:-/bin/sh
